# frigate-telegram-control.service
[Unit]
# __APPDIR__ автоматически заменяется на реальный путь при `./manage.sh install`.
Description=Frigate Telegram Monitor - Mute Controller (кнопки паузы)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__APPDIR__
ExecStart=/bin/bash -c "source __APPDIR__/venv/bin/activate && python __APPDIR__/mute_controller.py"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
