# frigate-telegram@.service
[Unit]
# Templated service: %i is a group name from config.py (GROUPS).
# __APPDIR__ is replaced with the real install path by `./manage.sh install`.
# Start any group with: systemctl enable --now frigate-telegram@group3
Description=Frigate Telegram Monitor - group %i
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__APPDIR__
ExecStart=/bin/bash -c "source __APPDIR__/venv/bin/activate && python __APPDIR__/frigate_telegram_monitor.py %i"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
